# club/scripts/supervisor/club_nastya.conf
; Установка (один раз на сервере, от root):
;   sudo cp scripts/supervisor/club_nastya.conf /etc/supervisor/conf.d/club_nastya.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Управление: supervisorctl start club_nastya:club_nastya_bot
; Деплой останавливает/поднимает: club_nastya:club_nastya_bot (см. deploy_prod.sh)

[program:club_nastya_bot]
command=/home/appuser/club_nastya/venv/bin/python3 /home/appuser/club_nastya/main.py
directory=/home/appuser/club_nastya
user=appuser
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/club_nastya_bot.err.log
stdout_logfile=/var/log/supervisor/club_nastya_bot.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED=1,PYTHONPATH="/home/appuser/club_nastya",PATH="/home/appuser/club_nastya/venv/bin:%(ENV_PATH)s"

[group:club_nastya]
programs=club_nastya_bot
priority=999
