# club/scripts/supervisor/rag_read_api.conf
[program:rag_read_api]
command=/home/appuser/club/venv/bin/python3 /home/appuser/club/scripts/run_rag_read_api.py --env-file /home/appuser/club/.env
directory=/home/appuser/club
user=appuser
autostart=true
autorestart=true
stopwaitsecs=10
stderr_logfile=/var/log/supervisor/rag_read_api.err.log
stdout_logfile=/var/log/supervisor/rag_read_api.out.log

; Доступ снаружи — только через reverse proxy / SSH tunnel.
; По умолчанию RAG_READ_API_HOST=127.0.0.1, порт 8790.
